# etc/supervisord.d/kibana4_restore.ini
[program:kibana4_restore]
command = /opt/qnib/kibana4/bin/restore.sh
redirect_stderr=true
startsecs=1
stdout_logfile = /var/log/supervisor/%(program_name)s.log
