# etc/supervisord.d/kibana4.ini
[program:kibana4]
command = /opt/qnib/kibana4/bin/start.sh
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
