# etc/supervisord.d/kibana.ini
[program:kibana]
command = /opt/qnib/kibana/bin/start_kibana.sh
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
