# etc/supervisord.d/kibana_restore.ini
[program:kibana_restore]
command = /opt/qnib/bin/kibana_restore.sh
redirect_stderr=true
startsecs=1
stdout_logfile = /var/log/supervisor/%(program_name)s.log
